# data/supervisor/treneris.conf
[program:accueil]
command=docker-compose -f docker-compose.yml up
autostart=true
autorestart=true
user=trener
directory=/home/trener/accueil/
