# data/supervisor/treneris.conf
[program:treneris-accueil]
command=docker-compose -d up
autorestart=true
user=trener
directory=/home/trener/accueil/
